# deploy/exe-dev/intent-heartbeat.timer
[Unit]
Description=Intent Computer Heartbeat Timer (every 15 minutes)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
